# kicker.socket
[Unit]
Description=django kicker socket

[Socket]
ListenStream=/home/chris/run/kicker.socket

[Install]
WantedBy=sockets.target
